# utilities/crontab
# Weather thrice a day
43 07,12,19 * * * /usr/local/bin/retrieve_weather.sh >/dev/null 2>&1
