# systemd/redis_sensor_bmx280.service
[Unit]
Description=Read from BMX280 and feed result into redis

[Service]
Type=simple
ExecStart=/usr/local/bin/redis_sensor.py

[Install]
WantedBy=multi-user.target
